# option/PackageDumpScript/firstboot.template
#!/bin/sh
#

# PROVIDE: firstboot
# REQUIRE: FILESYSTEMS
# KEYWORD: firstboot

name="firstboot"

. /etc/rc.subr

touch /firstboot-reboot

echo "Execute pkg pre/post/install script"

